# .config/autostart/alsarestore.desktop
[Desktop Entry]
Name=alsarestore
Exec=alsactl --file /home/vinz/.config/asound.state restore
Type=Application
