# .config/autostart/xset-keyboard.desktop
[Desktop Entry]
Name=Change keyboard keys repeat rate
Exec=xset r rate 170 25
Type=Application
