# systemd/.config/systemd/user/mouse-acceleration.service
[Unit]
Description=Set mouse acceleration
After=graphical-session.target

[Service]
ExecStart=/usr/bin/xset m 3 200
Restart=on-failure
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target

