# autostart/.config/autostart/wallpaper.desktop
#[Desktop Entry]
#Name=Set the wallpaper
#Exec=hsetroot -solid "#282828"
#Type=Application
