# autostart/.config/autostart/xset-escape.desktop
[Desktop Entry]
Name=Replace capslock by escape
Exec=setxkbmap -option caps:escape
Type=Application
