# autostart/.config/autostart/xset-mouse.desktop
[Desktop Entry]
Name=Enable mouse acceleration
Exec=xset m 3 200
Type=Application
